# layer/themes/infra_scripts/scripts/setup/demo/mate-dealer/mate-dealer-demo.service
[Unit]
Description=AI Dev Portal Demo - mate-dealer
After=network.target postgresql.service

[Service]
Type=simple
User=ormastes
WorkingDirectory=/home/ormastes/dev/aidev/scripts/setup/demo/mate-dealer
Environment="NODE_ENV=demo"
Environment="PORT=3303"
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
